# etc/systemd/systemd/hdd.service
[Unit]
Description=Suspend HDD

[Service]
ExecStart=/usr/bin/hdparm -q -S 120 -y /dev/sda
Type=idle

[Install]
WantedBy=multi-user.target
