# etc/systemd/systemd/undervolt.service
[Unit]
Description=Undervolt CPU

[Service]
ExecStart=/usr/bin/intel-undervolt apply
Type=idle

[Install]
WantedBy=multi-user.target
