# etc/systemd/systemd/getty@tty1.service.d/override.conf
[Service]
ExecStart=
ExecStart=-/usr/bin/agetty  --nonewline --noissue --autologin missu --noclear %I $TERM --noclear %I $TERM
