# etc/systemd/system/disable-usb-wakeup.service
[Unit]
Description=Disable USB Wakeup

[Service]
Type=oneshot
ExecStart=echo EHC1 > /proc/acpi/wakeup

[Install]
WantedBy=multi-user.target
